# robot/pi-runtime/services/vexy-camera.service
[Unit]
Description=Vexy System 2 camera broker

[Service]
Type=simple
WorkingDirectory=%h/llm-self-model-capstone/robot/pi-runtime
EnvironmentFile=%h/llm-self-model-capstone/robot/pi-runtime/config/defaults
EnvironmentFile=-%h/.config/vexy-system2/local
Environment=PYTHONPATH=%h/llm-self-model-capstone/robot/pi-runtime/src
ExecStart=/usr/bin/python3 -m vexy_system2.camera_broker
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
